# supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=npm start

[program:nginx]
command=nginx
